# crontab
#*/2 * * * * cd /data && node quittance.js >> /var/log/cron > ${STDOUT_LOC} 2> ${STDERR_LOC}
0 8 5 * * cd /data && /usr/local/bin/node quittance.js >> /var/log/cron

